# systemd/system/bird-update-roa.timer
[Unit]
Description=Update DN42 ROA Tables

[Timer]
OnBootSec=15min
OnUnitActiveSec=30min
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
