# systemd/system/bird-peer-finder.timer
[Unit]
Description=DN42 Peer Finder

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
